# docs/deployment-example/routeragent.service
[Unit]
Description=WiFi Manager router agent
After=network-online.target wifimanager.service dnsmasq-wifimanager.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/wifimanager-routeragent
EnvironmentFile=/etc/wifimanager/routeragent.env
ExecStart=/opt/wifimanager-routeragent/routeragent
Restart=on-failure
RestartSec=5s
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true

[Install]
WantedBy=multi-user.target
